# deploy/merylmoe.service
# To deploy manually:
#   scp deploy/merylmoe.service vps:/etc/systemd/system/merylmoe.service
#   ssh vps "systemctl daemon-reload && systemctl restart merylmoe"

[Unit]
Description=meryl.moe web server
After=network.target

[Service]
Type=simple
User=kuma
WorkingDirectory=/opt/meryl.moe
Environment=HOST=127.0.0.1
Environment=PORT=3000
Environment=LOG_DIR=/opt/meryl.moe/logs
ExecStart=/opt/meryl.moe/server
Restart=always

[Install]
WantedBy=multi-user.target
